# scripts/wifi-powersave-off.service
[Unit]
Description=Disable WiFi power management
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Prefer the modern `iw` tool; fall back to the deprecated `iwconfig` if `iw`
# can't set it (some older drivers only respond to iwconfig).
ExecStart=/bin/sh -c '/usr/sbin/iw dev WIFI_INTERFACE_PLACEHOLDER set power_save off || /usr/sbin/iwconfig WIFI_INTERFACE_PLACEHOLDER power off'

[Install]
WantedBy=multi-user.target
